# crontab
# Meridian Capital Partners — supercronic schedule (TZ=America/New_York)
# min hour dom mon dow   command
#
# Daily scoring: weekdays 17:15 ET. Refreshes prices, short interest, analyst
# estimates, earnings calendar, then re-scores all factors (~10-40 min).
15 17 * * 1-5 cd /app && python3 run_scoring.py --no-filings --no-13f >> /app/output/logs/scoring.log 2>&1
