# crontab
# Mediajedi Hedge Fund — supercronic schedule (TZ=America/New_York)
# min hour dom mon dow   command
#
# Daily scoring: weekdays 17:15 ET. Refreshes prices, short interest, analyst
# estimates, earnings calendar, then re-scores all factors (~10-40 min).
15 17 * * 1-5 cd /app && python3 run_scoring.py --no-filings --no-13f >> /app/output/logs/scoring.log 2>&1
#
# Intraday risk monitor: every 5 min during market hours (ET). Records NAV and fires
# circuit breakers (size-down/close-all/kill-switch/force-close). PROTECTIVE ONLY —
# it never opens new positions; portfolio execution stays manual (run_execution.py).
*/5 9-16 * * 1-5 cd /app && python3 -m execution.monitor >> /app/output/logs/monitor.log 2>&1
#
# Weekly full refresh: Sunday 18:00 ET (~2h). Deep data (filings/insider/13F) + Claude
# analysis (PAID, capped) + beta-neutral portfolio rebuild, so the book is fresh Monday
# morning. NEVER trades — execution is handled by the Monday job below.
0 18 * * 0 /app/scripts/weekly_refresh.sh
#
# Monday auto-execution: 09:45 ET (after the open). Runs the Sunday-built target book
# through the executor — veto- and kill-switch-gated, sized off LIVE equity — and emails
# a confirmation (or the reason it declined). This is the ONLY job that OPENS positions.
45 9 * * 1 cd /app && python3 -m scripts.auto_execute >> /app/output/logs/auto_execute.log 2>&1
#
# Daily report email: weekdays 17:00 ET (after close). Styled HTML recap — balance,
# today's P&L, book exposure, top movers. Read-only; never trades.
0 17 * * 1-5 cd /app && python3 -m scripts.daily_report >> /app/output/logs/daily_report.log 2>&1
